# image/proteos-guestagent.service
[Unit]
Description=ProteOS guest agent (PTY broker over vsock)
Documentation=https://github.com/tavon-ai/proteos
# Start once the vsock device and base filesystems are ready. The agent does
# not need the network (it serves over vsock), so we do not order on it.
After=local-fs.target
DefaultDependencies=yes

[Service]
Type=simple
ExecStart=/usr/local/bin/guestagent
# vsock:1024 + a login bash shell as root (Phase 3 decision #6). The control
# plane reaches this only through the host-side jailed uds — see decision #10
# (the trust boundary is the per-VM uds; per-machine app-layer identity arrives
# in Phase 5).
Environment=PROTEOS_GUEST_LISTEN=vsock:1024
Environment=PROTEOS_GUEST_SHELL=/bin/bash
Environment=PROTEOS_GUEST_SCROLLBACK_KIB=256
# Phase 4 (decision #7): persist mode. With PROTEOS_GUEST_PERSIST UNSET, the
# agent runs in DISK mode — it waits for the persistent disk, fsck's it, mounts
# it at /persist, and bind-mounts /persist/home→/root and /persist/workspace→
# /workspace so $HOME + the workspace live on the encrypted disk. The device is
# the node-agent-attached /dev/vdb. A missing disk degrades to ephemeral rather
# than refusing terminals. The agent runs as root (no User= here) so it can
# mount and, on resume, clock_settime + reseed the CRNG.
Environment=PROTEOS_GUEST_PERSIST_DEV=/dev/vdb
# Phase 8 (decision #4/#5): the code-server web forward. The agent listens on a
# second vsock port (1025) that the node-agent tunnel reaches on the guest web
# port, and lazily starts + supervises code-server (bound loopback-only on
# 127.0.0.1:13337) on the first web connection. code-server is baked at
# /usr/local/bin/code-server by image/build-rootfs.sh; a build that opted out
# (--no-codeserver) leaves this listener up but unable to start the editor.
Environment=PROTEOS_GUEST_WEB_LISTEN=vsock:1025
Environment=PROTEOS_CODESERVER_BIN=/usr/local/bin/code-server
# PP1 (port preview): the generic loopback forward. The agent listens on a third
# vsock port (1026) that the node-agent tunnel reaches on the guest preview port;
# for each connection the node-agent names a loopback port in a one-line preamble
# and the agent bridges to 127.0.0.1:<port>. No backend/supervisor — the user's
# own dev server is the backend; the control-plane gateway is the authenticator.
Environment=PROTEOS_GUEST_PREVIEW_LISTEN=vsock:1026
Restart=always
RestartSec=1
# The shell sessions are children of the agent; give them room.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
